# var/spool/cron/crontabs/root
0 0,12 * * * /usr/bin/php /home/epg2xml/epg2xml-web.php
